# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[program:ttslo-monitor]
command=/root/.local/bin/uv run python ttslo.py --interval 60 --config %(ENV_TTSLO_CONFIG_FILE)s --state %(ENV_TTSLO_STATE_FILE)s --log %(ENV_TTSLO_LOG_FILE)s
directory=/app
autostart=true
autorestart=true
startretries=3
stderr_logfile=/var/log/supervisor/ttslo-monitor.err.log
stdout_logfile=/var/log/supervisor/ttslo-monitor.out.log
environment=PATH="/root/.local/bin:%(ENV_PATH)s",KRAKEN_API_KEY="%(ENV_KRAKEN_API_KEY)s",KRAKEN_API_SECRET="%(ENV_KRAKEN_API_SECRET)s",KRAKEN_API_KEY_RW="%(ENV_KRAKEN_API_KEY_RW)s",KRAKEN_API_SECRET_RW="%(ENV_KRAKEN_API_SECRET_RW)s",TELEGRAM_BOT_TOKEN="%(ENV_TELEGRAM_BOT_TOKEN)s",TTSLO_CONFIG_FILE="%(ENV_TTSLO_CONFIG_FILE)s",TTSLO_STATE_FILE="%(ENV_TTSLO_STATE_FILE)s",TTSLO_LOG_FILE="%(ENV_TTSLO_LOG_FILE)s"

[program:ttslo-dashboard]
command=/root/.local/bin/uv run python dashboard.py --host 0.0.0.0 --port %(ENV_DASHBOARD_PORT)s
directory=/app
autostart=true
autorestart=true
startretries=3
stderr_logfile=/var/log/supervisor/ttslo-dashboard.err.log
stdout_logfile=/var/log/supervisor/ttslo-dashboard.out.log
environment=PATH="/root/.local/bin:%(ENV_PATH)s",KRAKEN_API_KEY="%(ENV_KRAKEN_API_KEY)s",KRAKEN_API_SECRET="%(ENV_KRAKEN_API_SECRET)s",KRAKEN_API_KEY_RW="%(ENV_KRAKEN_API_KEY_RW)s",KRAKEN_API_SECRET_RW="%(ENV_KRAKEN_API_SECRET_RW)s",TELEGRAM_BOT_TOKEN="%(ENV_TELEGRAM_BOT_TOKEN)s",DASHBOARD_PORT="%(ENV_DASHBOARD_PORT)s",TTSLO_CONFIG_FILE="%(ENV_TTSLO_CONFIG_FILE)s",TTSLO_STATE_FILE="%(ENV_TTSLO_STATE_FILE)s",TTSLO_LOG_FILE="%(ENV_TTSLO_LOG_FILE)s",TTSLO_CHECK_INTERVAL="60"
