# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[program:ttslo-monitor]
# Config paths are hardcoded to /config (Docker volume mount point)
# This matches the default CONFIG_DIR in entrypoint.sh
# Supervisord doesn't support environment variable interpolation in command strings,
# so we use fixed paths that match the Docker volume configuration (see docker-compose.yml volumes section)
# Environment variables (API keys, etc.) are inherited from parent process
command=/root/.local/bin/uv run python ttslo.py --interval 60 --config /config/config.csv --state /config/state.csv --log /config/logs.csv
directory=/app
autostart=true
autorestart=true
startretries=3
stderr_logfile=/var/log/supervisor/ttslo-monitor.err.log
stdout_logfile=/var/log/supervisor/ttslo-monitor.out.log

[program:ttslo-dashboard]
# Uses wrapper script to validate and inject DASHBOARD_PORT
command=/app/docker/run-dashboard.sh
directory=/app
autostart=true
autorestart=true
startretries=3
stderr_logfile=/var/log/supervisor/ttslo-dashboard.err.log
stdout_logfile=/var/log/supervisor/ttslo-dashboard.out.log
